# df-image-sdl/s6/xvfb/run
#!/bin/sh
exec Xvfb :0 -screen 0 ${XVFB_RESOLUTION:-1280x800x24} -nolisten tcp 2>&1
