# df-image-sdl/s6/df/run
#!/bin/sh
# Wait for Xvfb socket and VNC port before launching DF.
until [ -S /tmp/.X11-unix/X0 ]; do sleep 0.2; done
until nc -z 127.0.0.1 5900 2>/dev/null; do sleep 0.3; done

export DISPLAY=:0
cd /opt/df
exec ./dwarfort 2>&1
