# df-image-sdl/s6/df/run
#!/bin/sh
# Wait for Xvfb socket and VNC port before launching DF.
until [ -S /tmp/.X11-unix/X0 ]; do sleep 0.2; done
until nc -z 127.0.0.1 5900 2>/dev/null; do sleep 0.3; done
until [ -S /tmp/pulse/native ]; do sleep 0.2; done

export DISPLAY=:0
export PULSE_SERVER=unix:/tmp/pulse/native
export SDL_AUDIODRIVER=pulse
cd /opt/df

# Apply user-uploaded tilesets and selection. Non-fatal on error.
/apply-tilesets.sh || true

# Redirect gamelog.txt into the bind-mounted save dir for host-side tailing.
/link-gamelog.sh || true

# On SIGTERM (sent by s6 during container shutdown) forward the signal to DF
# so it exits promptly. The wrapper does NOT attempt to trigger an in-game
# save — saving is the player's responsibility and any blind keystroke
# injection would only work on the fortress play screen.
_term() {
    kill -TERM "$DF_PID" 2>/dev/null
    wait "$DF_PID" 2>/dev/null
    exit 0
}
trap _term TERM

# Launch through DFHack if installed, otherwise run DF directly.
if [ -x /opt/df/dfhack ]; then
    ./dfhack &
else
    ./dwarfort &
fi
DF_PID=$!

# Without a window manager, DF's SDL2 window never receives X11 input focus,
# so all keyboard and mouse button events from the VNC server are lost.
# Wait for the window to appear and give it focus explicitly.
xdotool search --sync --name "Dwarf Fortress" windowfocus &

wait "$DF_PID"
