# df-image-sdl/s6/vnc/run
#!/bin/sh
# Wait for Xvfb to create its socket before starting VNC.
until [ -S /tmp/.X11-unix/X0 ]; do sleep 0.2; done

# Force x11vnc to send a static arrow as the cursor shape to every client.
# Without -cursor arrow, x11vnc tries to read the real X cursor via XFIXES; on
# headless Xvfb with no window manager and no cursor font it falls back to its
# built-in "X" glyph, which is what the user sees in the browser. -cursor arrow
# bypasses that path entirely and is independent of xsetroot/xfonts-base.
exec x11vnc \
    -display :0 \
    -forever \
    -shared \
    -nopw \
    -rfbport 5900 \
    -noxdamage \
    -noipv6 \
    -xkb \
    -cursor arrow \
    -quiet 2>&1
