# df-image-sdl/s6/websockify/run
#!/bin/sh
# Wait for VNC to be listening.
until nc -z 127.0.0.1 5900 2>/dev/null; do sleep 0.3; done
exec websockify \
    --web /novnc \
    0.0.0.0:6080 \
    127.0.0.1:5900 2>&1
